# cron.weekly/maid_clean
0 13 * * 1 root maid clean -f
# Every Monday at 1PM
